# .docker/web/etc/supervisor/conf.d/laravel-queue.conf
[program:laravelqueue]
command=php /var/www/ukcp/artisan queue:listen
numprocs=1
autostart=true
autorestart=true
user=www
