# .docker/web/etc/supervisor/conf.d/laravel-queue-worker.conf
[program:laravelqueueworker]
command=php /var/www/ukcp/artisan queue:listen
numprocs=2
autostart=true
autorestart=true
user=www
process_name=%(program_name)s_%(process_num)02d
