# .docker/web/etc/supervisor/conf.d/laravel-queue-supervisor.conf
[program:laravelqueuesupervisor]
command=php /var/www/ukcp/artisan horizon
numprocs=1
autostart=true
autorestart=true
user=www
