# .docker/web/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=php-fpm
numprocs=1
autostart=true
autorestart=true
user=www
