# .docker/php/etc/supervisor/conf.d/websockets.conf
[program:websockets]
command=php /var/www/ukcp/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=www
